# tools/complexity_daemon/systemd/complexity-daemon.service
[Unit]
Description=Complexity Daemon
After=network.target

[Service]
ExecStart=%h/.local/bin/cogload daemon
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
